# install_scripts/deb/prerm
systemctl stop notrust-track.service
systemctl disable notrust-track.service
deluser --system notrust
delgroup --system notrust
