# resources/deb/prerm
systemctl stop zerotrust-track.service
systemctl disable zerotrust-track.service
deluser --system zerotrust
delgroup --system zerotrust
